# lib/support/upstart/redirus-worker.conf
start on starting redirus
stop on stopping redirus

