# lib/support/upstart/redirus.conf
pre-start script
  mkdir -p /home/atmosphere/proxyconf/redirus-worker/log
  chown -R atmosphere /home/atmosphere/proxyconf/redirus-worker/log

  mkdir -p /home/atmosphere/proxyconf/configs
  chown -R atmosphere /home/atmosphere/proxyconf/configs
end script
start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

